# deploy/driving-dashboard-sync.timer
[Unit]
Description=Daily Wialon driving safety sync timer

[Timer]
OnCalendar=*-*-* 06:10:00
Persistent=true
Unit=driving-dashboard-sync.service

[Install]
WantedBy=timers.target
